# os/services/menu.service
[Unit]
Description=Moab Menu
DefaultDependencies=false
After=docker.service

[Service]
Environment="MOABIAN=3.0.20"
Environment="PYTHONUNBUFFERED=1"
WorkingDirectory=/home/pi/moab/sw
ExecStart=/usr/bin/python3 menu.py --no-debug
Restart=no
User=pi

[Install]
WantedBy=multi-user.target
