# os/services/menu.service
[Unit]
Description=Moab Menu
DefaultDependencies=false
After=docker.service

[Service]
Type=forking
Environment="MOABIAN=3.0.19"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/moab/sw/menu.py --no-debug
WorkingDirectory=/home/pi/moab/sw
Restart=no
User=pi

[Install]
WantedBy=multi-user.target
